# mini-pc-setup/configs/gaming-stream-downscale.service
[Unit]
Description=Gaming WebRTC Stream (1440p → 1080p Downscale)
Requires=broadcast-box.service
After=broadcast-box.service

[Service]
Type=simple
User=streamer
WorkingDirectory=/home/streamer/streaming
Environment="STREAM_KEY=gaming"
Environment="SERVER_URL=http://localhost:8080/api/whip"
Environment="VIDEO_DEVICE=/dev/video0"
ExecStart=/home/streamer/streaming/stream-1080p-downscale.sh
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/streamer/streaming

[Install]
WantedBy=multi-user.target
